# contrib/systemd/fullbogons-nftables-gen.timer
[Unit]
Description=Run fullbogons-nftables-gen every 4 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
